# tilestream.conf
[program:tilestream]
user=root
command=/run_tilestream.sh
autorestart=true
stopsignal=INT
